# ff-api-update.cron
*/5 *   * * *   www-data    /usr/bin/php5 /opt/ff-api-update.php &>/dev/null
